# bin/consumer/consumer.service
[Unit]
Description=EPG service
After=network.target

[Service]
Type=simple
Restart=always
RestartSec=10s
ExecStart=/var/www/builds/consumer --listener 2 --config-file /var/www/app_metrica/config-consumer.yaml

[Install]
WantedBy=multi-user.target
